# scripts/caddy-docker.service
# SPDX-License-Identifier: AGPL-3.0-or-later
# Copyright (C) 2026 Relational Network
#
# Systemd service for Caddy Reverse Proxy (Docker)
#
# Installation:
#   sudo cp caddy-docker.service /etc/systemd/system/caddy-docker.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable caddy-docker
#   sudo systemctl start caddy-docker
#
# Logs:
#   sudo journalctl -u caddy-docker -f
#
# Prerequisites:
#   - Docker installed and running
#   - Ports 80 and 443 available (stop native Caddy if running)
#   - Environment file at /etc/caddy/environment with:
#       PILOT_DOMAIN=iob-staging.duckdns.org
#       DASHBOARD_ORIGIN=https://iob-dashboard.vercel.app
#
# Note: This uses Docker instead of native Caddy to match the CI/CD pipeline
# and ensure consistent behavior between local and deployed environments.

[Unit]
Description=Caddy Reverse Proxy (Docker)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=5
StartLimitInterval=300
StartLimitBurst=5
EnvironmentFile=/etc/caddy/environment

# Stop and remove any existing container
ExecStartPre=-/usr/bin/docker stop caddy
ExecStartPre=-/usr/bin/docker rm caddy

# Pull latest image with retry, fallback to cached image if pull fails
ExecStartPre=-/bin/bash -c 'docker pull ghcr.io/relational-network/relational-proxy:staging-latest || echo "Pull failed, using cached image"'

# Run the Caddy container
# - Binds ports 80 and 443 for HTTP/HTTPS
# - Mounts Caddy data/config for Let's Encrypt certs
# - Uses host networking for backend services
ExecStart=/usr/bin/docker run --rm --name caddy \
    --network host \
    -v caddy_data:/data \
    -v caddy_config:/config \
    -e PILOT_DOMAIN=${PILOT_DOMAIN} \
    -e DASHBOARD_ORIGIN=${DASHBOARD_ORIGIN} \
    ghcr.io/relational-network/relational-proxy:staging-latest

ExecStop=/usr/bin/docker stop caddy

# Health check - verify Caddy is listening on ports 80 and 443
ExecStartPost=/bin/bash -c 'for i in {1..30}; do sleep 1; ss -tln | grep -E ":(80|443) " && exit 0; done; echo "Caddy ports not listening after 30s" && exit 1'

[Install]
WantedBy=multi-user.target
